# linux_installer/qssweb2board_2.0-template/etc/profile.d/10-qssweb2board.sh
#!/bin/sh
export PATH=$PATH:/opt/QSSWeb2Board/
export QSSWEB2BOARD_SKETCHES=/var/tmp/QSSWeb2Board/sketches/
export QSSWEB2BOARD_BUILD=/var/tmp/QSSWeb2Board/build/
export QSSWEB2BOARD_KNOWNBOARDS=/etc/QSSWeb2Board/knownboards.json
export QSSWEB2BOARD_ARDUINO=/opt/QSSWeb2Board/res/arduino/
