# linux_gui_installer/qssweb2board_2.0-template/etc/profile.d/10-qssweb2board.sh
#!/bin/sh
export PATH=$PATH:/opt/QSSWeb2Board/
export QSSWEB2BOARD_SKETCHES=/opt/QSSWeb2Board/tmp/sketches/
export QSSWEB2BOARD_BUILD=/opt/QSSWeb2Board/tmp/build/
export QSSWEB2BOARD_KNOWNBOARDS=/opt/QSSWeb2Board/res/knownboards.json
export QSSWEB2BOARD_ARDUINO=/opt/QSSWeb2Board/res/arduino/
export QSSWEB2BOARD_TMP=/opt/QSSWeb2Board/tmp/
